# os/linux/etc/profile.d/zzz-koopa.sh
#!/bin/sh

# koopa shell
# https://koopa.acidgenomics.com/
KOOPA_PREFIX='/usr/local/koopa'
# shellcheck source=/dev/null
. "${KOOPA_PREFIX}/activate"
