# os/linux/etc/profile.d/zzz-koopa.sh
#!/bin/sh

# koopa shell
# https://koopa.acidgenomics.com/
# shellcheck source=/dev/null
. /usr/local/koopa/activate
