# config/etc/profile.d/koopa.sh
#!/bin/sh

# koopa shell
# https://github.com/acidgenomics/koopa
# shellcheck source=/dev/null
. /opt/koopa/activate
