# host/azure/config/etc/cron.d/cpi-reboot
@reboot /usr/local/koopa/host/azure/config/etc/reboot.sh
