# config/os/fedora/etc/profile.d/koopa.sh
# koopa shell
# https://github.com/acidgenomics/koopa
if [ -z "${KOOPA_SHELL:-}" ]
then
    # shellcheck source=/dev/null
    . /usr/local/koopa/activate
fi
