# host/aws/etc/profile.d/constellation.sh
#!/bin/sh

# """
# Constellation Pharma AWS VM shared shell configuration.
# Updated 2020-03-27 by Michael Steinbaugh.
# """

export KOOPA_CONFIG="constellation-aws"
